# etc/systemd/system/backhaul-iran100.service
[Unit]
Description=Backhaul Iran Port 100 (Iran)
After=network.target

[Service]
Type=simple
ExecStart=/root/backhaul-core/backhaul_premium -c /root/backhaul-core/iran100.toml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
